# dnsproxy.service
[Unit]
Description=AdGuard DNS-over-QUIC Proxy
After=network-online.target
Before=pihole-FTL.service
[Service]
Type=simple
# Usiamo lhostname per lidentificazione corretta con Quad9
ExecStart=/usr/local/bin/dnsproxy \
  -l 127.0.0.1 -p 5353 \
  -u quic://9.9.9.11:853 \
  -u quic://149.112.112.11:853 \
  --pending-requests-enabled \
  -r 500 --cache --cache-optimistic
Restart=always
RestartSec=5
User=root
[Install]
WantedBy=multi-user.target
